# rootfs/etc/s6/docker/run
#!/bin/bash
exec docker daemon --insecure-registry 10.0.3.4:5000